# dockerfiles/operator/operator.conf
[program:operator]
command=/app/run_operator.sh
autostart=true
autorestart=true
stderr_logfile=/var/log/operator.err.log
